# viz/script/vizinit.service
[Unit]
Description=Arcviz server

[Service]
ExecStart=/etc/init.d/vizinit start

[Install]
WantedBy=multi-user.target